# systemd/raspberrymeet.service
[Unit]
Description=RaspberryMeet Orchestrator Service
Documentation=https://github.com/Sico93/RaspberryMeet
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/RaspberryMeet

# Environment
Environment="PATH=/home/pi/RaspberryMeet/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/pi/RaspberryMeet/.env

# Start the orchestrator
ExecStart=/home/pi/RaspberryMeet/venv/bin/python -m src.orchestrator.main

# Restart policy
Restart=always
RestartSec=10

# Resource limits
MemoryMax=1G
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raspberrymeet

# Security
NoNewPrivileges=true
PrivateTmp=true

# GPIO access (add user to gpio group instead of running as root)
SupplementaryGroups=gpio video

[Install]
WantedBy=multi-user.target
